# crontab
*/5 *       * * *   raptors ~/media-scripts/music_sort/music_sort.py
1-60/5 *    * * *   raptors ~/media-scripts/music_sort/playlist_import.py
